# util/systemd/weewx.service
# systemd configuration for weewx

[Unit]
Description=weewx weather system
Requires=syslog.service
Requires=ntp.service
After=syslog.service
After=ntp.service
RequiresMountsFor=/home

[Service]
ExecStart=/home/weewx/bin/weewxd --daemon --pidfile=/var/run/weewx.pid /home/weewx/weewx.conf
ExecReload=/bin/kill -HUP $MAINPID
Type=simple
PIDFile=/var/run/weewx.pid
#User=weewx
#Group=weewx

[Install]
WantedBy=multi-user.target
